# security/systemd/disk-check.service
[Unit]
Description=Disk usage alert

[Service]
Type=oneshot
EnvironmentFile=/etc/server.env
ExecStart=/usr/local/sbin/disk-check
